# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:apollod]
command=/usr/local/bin/apollod
